# worker/lumina-push.timer
# 每 30 分钟触发一次 worker（once）。机器醒着即按时推送；睡眠期错过的，worker 内 catch-up 会补。
[Unit]
Description=Run Lumina Feed push worker periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
