# worker/lumina-push.service
# systemd 服务（once 模式，由 .timer 触发）。放 ~/.config/systemd/user/ 后:
#   systemctl --user daemon-reload && systemctl --user enable --now lumina-push.timer
[Unit]
Description=Lumina Feed push worker (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/lumina-feed
Environment=LUMINA_CONFIG=%h/lumina-feed/worker/config.json
# 密钥经 systemd 环境文件注入（chmod 600），不入仓库：
EnvironmentFile=%h/.config/lumina/secrets.env
ExecStart=/usr/bin/node --experimental-strip-types %h/lumina-feed/worker/push-worker.mjs --once
